# root/etc/services.d/cron/run
#!/usr/bin/with-contenv sh
# Regular daily/weekly/monthly cronjobs to run
/usr/sbin/crond -f -S -l 0 -c /etc/crontabs
